# update-points-and-ranks.service
[Service]
ExecStart=%h/unique-admin/update-points-and-ranks
