# update-points-and-ranks.timer
[Timer]
OnBootSec=1min
OnUnitInactiveSec=10min

[Install]
WantedBy=timers.target
