# update-and-announce-records.service
[Service]
ExecStart=%h/unique-admin/update-and-announce-records
