# website.service
[Service]
WorkingDirectory=%h/unique-web
ExecStart=node bin/www.js

[Install]
WantedBy=default.target
